# deploy/systemd/china-succession-notices.timer
[Unit]
Description=Run China Succession notice sync every 30 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=china-succession-notices.service

[Install]
WantedBy=timers.target
